# contrib/ip_detector.openrc
#!/sbin/openrc-run

name="IP Detector"
description="Public IP monitoring with Telegram notifications"
command="/usr/local/bin/ip_detector"
command_args="--daemon --interval 300"
command_user="${RC_SVCNAME##*.}"
command_background="yes"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    need net
    after firewall
}
